# ops/com.corelifecycle.catchup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.corelifecycle.catchup</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/jacobpazhoor/CoreLifecycle/scripts/daily_autopilot.sh</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Library/Frameworks/Python.framework/Versions/3.12/bin:/usr/local/bin:/Users/jacobpazhoor/.local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <!-- CATCH-UP: if the 2am run missed/failed (Mac asleep, render killed, etc.), retry at boot/login
       and every 2h. The runner self-guards via runs/last_post.txt — it does NOTHING if a video
       already went out today, so this is safe to fire often. -->
  <key>RunAtLoad</key><true/>
  <key>StartInterval</key><integer>7200</integer>
  <key>StandardOutPath</key><string>/Users/jacobpazhoor/CoreLifecycle/runs/autopilot/catchup.out</string>
  <key>StandardErrorPath</key><string>/Users/jacobpazhoor/CoreLifecycle/runs/autopilot/catchup.err</string>
  <key>ProcessType</key><string>Background</string>
  <key>LowPriorityIO</key><true/>
</dict>
</plist>
